# ndoutils-1.4b7/debian/ndoutils-NG-DB.ndoutils.init.in
#! /bin/sh

### BEGIN INIT INFO
# Provides:          ndoutils-@@NGV@@
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Should-Start:      $local_fs
# Should-Stop:       $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts ndo2db daemon for @@NG@@ at boot time
# Description:       This will start the ndo2db daemon for @@NG@@ at boot time
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/sbin/ndo2db
NAME=ndoutils
DESC=ndoutils

test -x $DAEMON || exit 0

# Include ndoutils defaults if available
if [ -f /etc/default/ndoutils ] ; then
	. /etc/default/ndoutils
fi

set -e

case "$1" in
  start)
	echo -n "Starting $DESC: "
	test ! "$ENABLE_NDOUTILS" = 1 && { echo "Not enabled in /etc/default/ndoutils"; exit 0; }
	start-stop-daemon --start --quiet --user nagios --group nagios --exec $DAEMON --background -- $DAEMON_OPTS
	echo "$NAME."
	;;
  stop)
	echo -n "Stopping $DESC: "
	start-stop-daemon --stop --oknodo --quiet --exec $DAEMON 
	echo "$NAME."
	;;
  #reload)
	#
	#	If the daemon can reload its config files on the fly
	#	for example by sending it SIGHUP, do it here.
	#
	#	If the daemon responds to changes in its config file
	#	directly anyway, make this a do-nothing entry.
	#
	# echo "Reloading $DESC configuration files."
	# start-stop-daemon --stop --signal 1 --quiet --exec $DAEMON
  #;;
  force-reload)
	#
	#	If the "reload" option is implemented, move the "force-reload"
	#	option to the "reload" entry above. If not, "force-reload" is
	#	just the same as "restart" except that it does nothing if the
	#   daemon isn't already running.
	# check wether $DAEMON is running. If so, restart
	start-stop-daemon --stop --test --quiet --exec $DAEMON \
	&& $0 restart \
	|| exit 0
	;;
  restart)
    echo -n "Restarting $DESC: "
	start-stop-daemon --stop --oknodo --quiet --exec $DAEMON 
	sleep 1
	echo "$NAME."
	test ! "$ENABLE_NDOUTILS" = 1 && { echo "Not enabled in /etc/defaults/ndoutils"; exit 0; }
	start-stop-daemon --start --quiet --exec $DAEMON --background -- $DAEMON_OPTS
	echo "$NAME."
	;;
  *)
	N=/etc/init.d/$NAME
	# echo "Usage: $N {start|stop|restart|reload|force-reload}" >&2
	echo "Usage: $N {start|stop|restart|force-reload}" >&2
	exit 1
	;;
esac

exit 0
